# etc/rdma-ndd.service.in
[Unit]
Description=RDMA Node Description Daemon
After=rdma.service

[Service]
Type=forking
EnvironmentFile=@IBDIAG_CONFIG_PATH@/ibdiag.conf
ExecStart=@exec_prefix@/sbin/rdma-ndd

[Install]
WantedBy=multi-user.target

