# ibacm/ibacm.service.in
[Unit]
Description=Starts the InfiniBand Address Cache Manager daemon
Documentation=man:ibacm file:@CMAKE_INSTALL_SYSCONFDIR@/rdma/ibacm_opts.cfg
After=opensm.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/ibacm -P

[Install]
WantedBy=network.target
