# rdma-ndd/rdma-ndd.service.in
[Unit]
Description=RDMA Node Description Daemon
Documentation=man:rdma-ndd

[Service]
Type=notify
Restart=always
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/rdma-ndd --systemd

